# config/launchd/com.mac-messages-mcp.prod.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Production launchd service for Mac Messages MCP.
  Installed to ~/Library/LaunchAgents/ by scripts/setup-services.sh

  DO NOT edit this file directly – run scripts/setup-services.sh to
  regenerate it with the correct paths for your machine.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mac-messages-mcp.prod</string>

  <key>ProgramArguments</key>
  <array>
    <string>__UV_PATH__</string>
    <string>run</string>
    <string>python</string>
    <string>scripts/run_oauth_server.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_PATH__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PORT</key>
    <string>8000</string>
    <key>EXTERNAL_URL</key>
    <string>__PROD_EXTERNAL_URL__</string>
  </dict>

  <!-- Start automatically when the user logs in -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Restart automatically on crash -->
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__REPO_PATH__/logs/prod.log</string>

  <key>StandardErrorPath</key>
  <string>__REPO_PATH__/logs/prod.err</string>
</dict>
</plist>
