# deployment/imac-launchagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS LaunchAgent to auto-start the Mac Messages MCP remote server.

  Install:
    1. Edit the ProgramArguments below (set your --server-url)
    2. Copy to ~/Library/LaunchAgents/:
         cp imac-launchagent.plist ~/Library/LaunchAgents/com.macmessagesmcp.remote.plist
    3. Load:
         launchctl load ~/Library/LaunchAgents/com.macmessagesmcp.remote.plist
    4. Verify:
         launchctl list | grep macmessagesmcp

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.macmessagesmcp.remote.plist
    rm ~/Library/LaunchAgents/com.macmessagesmcp.remote.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.macmessagesmcp.remote</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Adjust path to your uv binary if needed -->
        <string>/usr/local/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>-m</string>
        <string>mac_messages_mcp.remote</string>
        <string>--server-url</string>
        <string>https://mcp.yourdomain.com</string>
        <string>--port</string>
        <string>8000</string>
    </array>

    <!-- Set the working directory to the repo -->
    <key>WorkingDirectory</key>
    <string>/Users/YOU/mac_messages_mcp</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/mac-messages-mcp-remote.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/mac-messages-mcp-remote.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
